# installer/templates/ploydok.target
[Unit]
Description=Ploydok
Requires=ploydok.service
After=ploydok.service

[Install]
WantedBy=multi-user.target
