# installer/templates/ploydok.service
[Unit]
Description=Ploydok control-plane (Docker Swarm stack)
Requires=docker.service ploydok-port-isolation.service
After=docker.service ploydok-port-isolation.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=${PLOYDOK_INSTALL_DIR}
EnvironmentFile=${PLOYDOK_DATA_DIR}/.env
# `docker stack deploy` is idempotent: re-running it after the systemd unit
# is restarted reconciles the stack to the current docker-stack.yml without
# tearing existing services down. The first invocation provisions the stack;
# subsequent ones (e.g. after a reboot) are essentially no-ops.
ExecStart=/usr/bin/docker stack deploy --resolve-image always --with-registry-auth --prune -c ${PLOYDOK_INSTALL_DIR}/docker-stack.yml ploydok
ExecStop=/usr/bin/docker stack rm ploydok
TimeoutStartSec=0

[Install]
WantedBy=ploydok.target
