# installer/templates/ploydok-update.service
[Unit]
Description=Ploydok control-plane image refresher
Requires=docker.service
After=docker.service ploydok.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=${PLOYDOK_INSTALL_DIR}
EnvironmentFile=${PLOYDOK_DATA_DIR}/.env
# The helper verifies every image against the pinned GitHub Actions identity
# before pulling or asking Swarm to resolve a moving tag such as :edge.
ExecStart=/usr/local/lib/ploydok/update-stack.sh
TimeoutStartSec=300
